# config/config/includes.chroot/usr/lib/systemd/user/pacific-linux-notify-updates.service
[Unit]
Description=Pacific Linux — notify about available updates

[Service]
Type=oneshot
ExecStart=/usr/local/bin/pacific-linux-notify-updates
