# config/config/includes.chroot/usr/lib/systemd/user/venu-pacific-notify-updates.timer
[Unit]
Description=Check for a Venu Pacific update notification periodically

[Timer]
OnStartupSec=2min
OnUnitActiveSec=1h
Persistent=true

[Install]
WantedBy=timers.target
