# config/config/includes.chroot/usr/lib/systemd/user/venu-pacific-llama-prewarm.service
[Unit]
Description=Venu Pacific AI assistant prompt prewarm
# Runs once at login, after the model server unit: waits for the server
# to finish loading, then evaluates the fixed system-prompt + tool-schema
# prefix (~3,500 tokens, ~2 minutes on slow hardware) into its KV cache
# in the background. Without this, the FIRST assistant reply after every
# boot paid that cost while the user watched "Thinking..." -- confirmed
# repeatedly in real sessions. With it, the cost is paid invisibly at
# login, and the first reply only processes the user's own message.
#
# Type=simple (not oneshot) so login/session startup is never blocked
# waiting the ~2 minutes this takes -- it just runs alongside.
After=venu-pacific-llama.service
Wants=venu-pacific-llama.service

[Service]
Type=simple
ExecStart=/usr/local/bin/venu-pacific-assistant --prewarm
Nice=10
IOSchedulingClass=idle

[Install]
WantedBy=default.target
