# config/config/includes.chroot/etc/systemd/system/pacific-linux-update-check.service
[Unit]
Description=Pacific Linux — check for available updates
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/pacific-linux-update-check
