# config/config/includes.chroot/etc/xdg/autostart/venu-pacific-welcome.desktop
[Desktop Entry]
Type=Application
Name=Welcome to Venu Pacific
Comment=First-run orientation (runs once)
Exec=/usr/local/bin/venu-pacific-welcome
Icon=venu-pacific
OnlyShowIn=XFCE;
NoDisplay=true
X-GNOME-Autostart-Phase=Application
